# supervisor/go-space-chat.ini
[program:go-space-chat]
directory=/www/go-space-chat/
command= go run main.go

startsecs=1                   ; 也是一样，进程启动后跑了几秒钟，才被认定为成功启动，默认1
startretries=5                ; 失败最大尝试次数，默认3

redirect_stderr=true
stdout_logfile=AUTO
